# graphanainstall/templates/grafana.service.j2
[Unit]
Description=Grafana
After=network.target

[Service]
User=grafana
Group=grafana
ExecStart=/usr/sbin/grafana-server web
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
